# addon/rootfs/etc/services.d/royaframe_bridge/finish
#\!/command/with-contenv sh
echo "[royaframe_bridge] stopped"
