# addon/rootfs/etc/services.d/royaframe_bridge/run
#\!/command/with-contenv sh
set -e
echo "[royaframe_bridge] starting..."
exec node /app/bridge/src/server.js
