# addon/rootfs/etc/services.d/royaframe_bridge/run
#\!/usr/bin/with-contenv bashio
# RoyaFrame Bridge service run script

bashio::log.info "Starting RoyaFrame Bridge..."

cd /app
exec node src/server.js
